# packaging/linux/desktop-agent-guardian.service
[Unit]
Description=Desktop Agent 守护进程（控制面：巡检与自愈编排）
# 依赖主服务先起来；主服务挂了由它自己的 unit 拉起，本单元只管守护进程自身
# 注意：desktop-agent.service 是主服务 unit（远端已存在，由它管理主服务）。
After=network.target desktop-agent.service
Wants=desktop-agent.service

[Service]
Type=simple
# 安装路径/参数从环境变量读取，未配置时使用下方默认值。
# 覆盖方式：把变量写进 /etc/default/desktop-agent-guardian（本 unit 自动加载）
#   echo 'DA_PYTHON=/srv/desktop-agent/.venv/bin/python'   > /etc/default/desktop-agent-guardian
#   echo 'SELF_HEAL_HEALTH_URL=http://127.0.0.1:8080/health' >> /etc/default/desktop-agent-guardian
EnvironmentFile=-/etc/default/desktop-agent-guardian
# 注意：WorkingDirectory 不支持 ${ENV} 展开，故省略；daemon 不依赖 cwd（路径靠 __file__ 推导）
ExecStart=${DA_PYTHON:-/opt/desktop-agent/.venv/bin/python} -m agent_core.guardian_daemon
# 恢复时执行的重启命令（值含空格必须加引号）；主服务用 systemd 管理：systemctl restart desktop-agent
Environment="SELF_HEAL_RESTART_CMD=systemctl restart desktop-agent"
# 健康探针 URL 默认 http://127.0.0.1:8899/health（见 guardian_daemon.py 默认值）。
# 若主 app 监听其它端口（如本机 8080），在 /etc/default/desktop-agent-guardian 写一行覆盖即可。
Restart=always
RestartSec=5
# 守护进程自身崩了由 systemd 拉起，确保"控制面"常驻
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
